# openiked-keygen.target
[Unit]
Wants=openiked-keygen.service
PartOf=openiked.service
